# supervisord.conf
; Single-container process manager.
; Runs the main API, the LoanPASS pricing service, the Vite frontend, and a
; pricing health watchdog as INDEPENDENT programs. Each has autorestart=true, so
; a crash restarts ONLY that program — the others keep running.
;
; Restart coverage for pricing:
;   - CRASH (process exits)  -> supervisor autorestart (below).
;   - HANG  (alive but stuck, e.g. LoanPASS rate-limit/token churn) -> the
;     `pricing-watchdog` program polls /api/loanpass/health and runs
;     `supervisorctl restart pricing`. Every action is logged ("PRICING
;     WATCHDOG" / "PRICING SERVICE STARTED"), so restarts are visible in logs.
;
; Manual restart of just pricing (no container bounce):
;   docker exec <container> supervisorctl restart pricing

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

; Enable `supervisorctl` inside the container (restart individual programs).
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:api]
command=/opt/venv/bin/python -m uvicorn backend.api:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pricing]
command=/opt/venv/bin/python -m uvicorn backend.pricing_app:app --host 0.0.0.0 --port 8090
directory=/app
autostart=true
autorestart=true
; Keep retrying if LoanPASS startup hiccups; supervisor backs off between tries.
startretries=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=npm run dev -- --host 0.0.0.0 --port 5173
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Polls the pricing /health endpoint and restarts ONLY the pricing program when
; it hangs (supervisor's autorestart can't see a hang — the process is alive).
[program:pricing-watchdog]
command=/opt/venv/bin/python -m backend.tools.pricing_watchdog
directory=/app
autostart=true
autorestart=true
; Don't start probing until pricing has had a moment to boot.
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
